# mcu_shutdown
#!/bin/sh /etc/rc.common

STOP=8

stop() {
        echo '{"system":"reboot"}' > /dev/ttyS0
        sleep 1
}